# systemd/dmx_audio_react.service
[Unit]
Description=DMX Audio Reactive Light Controller
After=oled_splash.service network-online.target sound.target olad.service
Wants=oled_splash.service network-online.target olad.service

[Service]
Type=simple
User=benglasser
WorkingDirectory=/home/benglasser/pi-dmx-controller-v2
Environment=PYTHONUNBUFFERED=1 ENABLE_TUI=0
# Skip startup if dev mode flag exists
ExecStartPre=/bin/sh -c '! test -f /tmp/dmx-dev-mode'
ExecStart=/home/benglasser/pi-dmx-controller-v2/.venv/bin/python /home/benglasser/pi-dmx-controller-v2/dmx_audio_react.py --no-tui
Restart=on-failure
RestartSec=3
AmbientCapabilities=CAP_SYS_NICE
Nice=-5

[Install]
WantedBy=multi-user.target