# systemd/oled_splash.service
[Unit]
Description=OLED splash (Loading...) before DMX app
# Run as soon as root is mounted - no oled_wake (initramfs handles early display)
After=local-fs.target
Before=multi-user.target

[Service]
Type=oneshot
# Run as same user as dmx (adjust User/WorkingDirectory to match your install)
User=pi
WorkingDirectory=/home/pi/pi-dmx-controller-v2
ExecStart=/home/pi/pi-dmx-controller-v2/.venv/bin/python /home/pi/pi-dmx-controller-v2/oled_boot.py
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
