# systemd/pi-dmx.service
[Unit]
Description=Pi DMX Controller
After=network-online.target olad.service oled_splash.service
Wants=olad.service

[Service]
Type=simple
# Set User and paths to match your install (default: pi user, ~/pi-dmx-controller-v2)
User=pi
WorkingDirectory=/home/pi/pi-dmx-controller-v2
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/pi/pi-dmx-controller-v2/.venv/bin/python /home/pi/pi-dmx-controller-v2/dmx_audio_react.py
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
